# ops/systemd/primal-backup.service
[Unit]
Description=Primal nightly database backup
After=postgresql.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/primal-backup
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
